# dockerfiles/files/billing-queries.service
[Unit]
Description=Billing queries update script

[Service]
Type=simple
ExecStart=/root/bin/billing-queries.js
